# usr/lib/systemd/system/mythical-pve.service
[Unit]
Description=3RD Party License Server

[Service]
Type=oneshot
ExecStart=/usr/bin/mythical-pve
